# docker/monitorq/etc/service/celery/run
#!/bin/bash -e

cd /code
exec chpst -u www-data:www-data /usr/local/bin/celery worker -A nodewatcher -l info -Q monitor -B 2>&1
