# docker/monitor/etc/service/monitor/run
#!/bin/bash -e

cd /code
exec chpst -u www-data:www-data /usr/bin/python manage.py monitord 2>&1
